# docker_files/supervisor/supervisor.conf
[supervisord]
nodaemon=true
logfile=docker_files/supervisor/supervisord.log

[program:php_fpm]
command=php-fpm
autostart=true
autorestart=true

[include]
files = /etc/supervisor/conf.d/*.conf

